# concurrent.futures_/satellite.ini
[program:satellite]
directory=repo
command=venv/bin/gunicorn -b 0.0.0.0:5000 -w 3 app:create_app()
user=david
autostart=false
autorestart=false
startretries=0
stdout_logfile=satellite.log
stderr_logfile=satellite.log
