# supervisor_/satellite3.ini
[program:satellite3]
command=python3 simple_flask_app.py
user=david
autostart=False
autorestart=False
